# crontab
*/30 * * * *	root	java -Duser.home=/data -jar /root/mytube.jar /var/www/html/ > /data/mytube.log 2>&1
